# roles/configure_md_raid/templates/mdcheck_template_timer.timer.j2
[Unit]
Description=Initiates array check on {{ item.device }}

[Timer]
OnCalendar={{ item.frequency }}
Persistent=true

[Install]
WantedBy=timers.target
Also=mdcheck_continue.timer